# zope/files/cron.weekly/zodb-pack
#!/bin/sh

export PYTHONPATH=/usr/lib/zope2.9/lib/python

for socket in /srv/zope2.9/zeo/*/var/zeo.sock; do
	/usr/lib/zope2.9/bin/zeopack.py -d 7 -U $socket
done
